# conf/supervisor/supervisord.conf
; Supervisor -- used by "devfast", for a deployment like setup on Joel's laptop
;
; start as sudo supervisord -c conf/supervisor/supervisord.conf

[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:9001        ; (ip_address:port specifier, *:port for all iface)

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://127.0.0.1:9001 ; use an http:// url to specify an inet socket

; The below sample program section shows all possible program subsection values,
; create one or more 'real' program: sections to be able to control them under
; supervisor.

[program:uwsgi]
command=/Users/joel/z/venv/bin/uwsgi -i /Users/joel/z/conf/uwsgi/devfast-joel.ini
stopsignal=INT

[program:nginx]
command=/opt/local/sbin/nginx -c /Users/joel/z/conf/nginx/devfast-joel.conf

[program:memcached]
command=/usr/bin/memcached -u joel

[program:rabbitmq]
;command=/opt/local/sbin/rabbitmq-server
command=/Users/joel/z/dev/scripts/rabbitmq.sh